# docker/crontab
* * * * * for i in $(seq 0 59); do curl -s http://localhost:8080/api/hello-world/hello >/dev/null 2>&1; sleep 1; done